# upstart/fitpc/uwesub-fitpc.conf
# uwesub ros upstart

# When to start the service
start on local-filesystems

#stop on

respawn

script
        exec /home/uwesub/projectChmaera/upstart/fitpc/fitpcupstart.sh
end script
