# upstart/roboard/uwesub-roboard.conf
# uwesub ros upstart

# When to start the service
start on local-filesystems

#stop on

respawn

script
        exec /home/uwesub/projectChmaera/upstart/roboard/roboardupstart.sh
end script
